# package/busybox/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/bash
::askfirst:-/bin/bash
::cttlaltdel:/bin/umount -a -r
